# configs/automatic-backup/checkmk-deployment-backup.service
[Unit]
Description=Checkmk Deployment Backup

[Service]
ExecStart=/usr/bin/checkmk-deployment-backup.sh
